# systemd/djconnect-nightly-reboot.service
[Unit]
Description=Nightly reboot for DJConnect Raspberry Pi

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reboot -i
