# systemd/djconnect-client.service
[Unit]
Description=DJConnect Raspberry Pi touch client
After=network-online.target systemd-user-sessions.service
Wants=network-online.target

[Service]
Type=simple
User=djconnect
Group=djconnect
WorkingDirectory=/opt/djconnect
Environment=HOME=/home/djconnect
Environment=XDG_RUNTIME_DIR=/run/user/%U
ExecStart=/usr/bin/xinit /opt/djconnect/current/bin/djconnect-pi-client --config /opt/djconnect/config/client.json -- :0 -nolisten tcp -nocursor
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
