# systemd/djconnect-client.service
[Unit]
Description=DJConnect Raspberry Pi touch client
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=djconnect
Group=djconnect
Environment=DISPLAY=:0
ExecStart=/opt/djconnect/current/bin/djconnect-pi-client --config /opt/djconnect/config/client.json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target

